# systemd/mm-trade-multi.service
[Unit]
Description=MM-Profit BotV4 Multi-Wallet Instance %i
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/path/to/MM-profit
Environment="NODE_ENV=production"
Environment="WALLET_INDEX=%i"

# Main bot command with wallet index
ExecStart=/usr/bin/node src/bots/botv4.js

# Restart configuration
Restart=always
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

# Logging (separate log per wallet)
StandardOutput=append:/var/log/mm-trade/botv4-%i.log
StandardError=append:/var/log/mm-trade/botv4-%i-error.log
SyslogIdentifier=mm-trade-v4-%i

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/path/to/MM-profit/data

# Resource limits
LimitNOFILE=65536
MemoryMax=500M

[Install]
WantedBy=multi-user.target